# installer/linux/sonde-gateway.service
[Unit]
Description=Sonde Gateway Service
Documentation=https://github.com/Alan-Jowett/sonde
After=network.target

[Service]
Type=simple
User=sonde
Group=sonde
SupplementaryGroups=dialout
WorkingDirectory=/var/lib/sonde
# The serial port must be configured in /etc/sonde/environment via SERIAL_PORT=<device>.
EnvironmentFile=/etc/sonde/environment
ExecStart=/usr/local/bin/sonde-gateway \
    --db /var/lib/sonde/gateway.db \
    --port ${SERIAL_PORT} \
    --key-provider secret-service
Restart=on-failure
RestartSec=5s
# Create /var/run/sonde/ for the admin socket (default: /var/run/sonde/admin.sock).
RuntimeDirectory=sonde
RuntimeDirectoryMode=0750

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/sonde
ReadOnlyPaths=/etc/sonde
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
